# garbage-collector/crontab
* * * * * docker exec registry registry garbage-collect /etc/docker/registry/config.yml --delete-untagged > /proc/1/fd/1 2>/proc/1/fd/2
